# templates/synapse/systemd.target.j2
[Unit]
Description=Synapse {{ instance.server_name }}.parent.target
After=network.target

[Install]
WantedBy=multi-user.target
